# .config/systemd/user/borgbackup.service
[Unit]
Description=Daily backup
ConditionACPower=true

[Service]
EnvironmentFile=%h/.bashrc_private
Type=oneshot
ExecStart=/usr/bin/borgmatic

